# flumeRest/supervisor.conf
# modify [program:flume-flask-service].directory
# modify [include].files

[unix_http_server]
file=/tmp/flume_services_supervisor.sock          ; path to your socket file

[supervisord]
logfile=%(here)s/supervisord.log                ; supervisord log file
logfile_maxbytes=100MB                          ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=debug                                  ; info, debug, warn, trace
pidfile=%(here)s/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=%(here)s                            ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/flume_services_supervisor.sock  ; use a unix:// URL  for a unix socket

[program:flume-flask-service]
command=python flumeAPI.py
directory=/path/to/service
user = centos
autorestart = true
startsecs = 3
stderr_logfile=/tmp/flume-flask-service.err
stdout_logfile=/tmp/flume-flask-service.out

[include]
#files=%(here)s/conf.d/*.conf                  ; not work
files=/path/to/service/conf.d/*.conf

